# supervisord.conf
[supervisord]
nodaemon=true

[program:squid]
command=/sbin/entrypoint.sh

[program:saboteur]
command=/usr/bin/saboteur-agent
